# cookbook/mongodb/templates/mongod.service.erb
[Unit]
Description=High-performance, schema-free document-oriented database
Documentation=man:mongod(1)
After=network.target

[Service]
Type=forking
User=mongodb
Group=mongodb
RuntimeDirectory=mongod
PIDFile=/var/run/mongod/mongod.pid
ExecStart=/usr/bin/mongod --quiet --config /etc/mongod.conf
KillMode=mixed

[Install]
WantedBy=multi-user.target
